# n1/files/etc/hotplug.d/iface/99-ipv6-ra
#!/bin/sh

if [ "$ACTION" = "ifup" ] && [ "$DEVICE" = "br-lan" ]; then
    echo 2 >/proc/sys/net/ipv6/conf/br-lan/accept_ra
fi
