# n1/files/etc/uci-defaults/99-cron_reboot
#!/bin/sh

TARGET="/etc/init.d/cron_reboot"

if [ -f "$TARGET" ]; then
    if [ ! -x "$TARGET" ]; then
        chmod +x "$TARGET"
    fi
    "$TARGET" enable
fi

exit 0
